# includes.container/theme-switcher/etc/skel/.config/systemd/user/theme-startup.service
[Unit]
Description=Set theme based on time of day
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/theme-switch auto

[Install]
WantedBy=graphical-session.target
