# includes.container/theme-switcher/etc/skel/.config/systemd/user/theme-day.service
[Unit]
Description=Switch to day (light) theme

[Service]
Type=oneshot
ExecStart=/usr/local/bin/theme-switch day
