# includes.container/theme-switcher/etc/skel/.config/systemd/user/theme-night.service
[Unit]
Description=Switch to night (dark) theme

[Service]
Type=oneshot
ExecStart=/usr/local/bin/theme-switch night
